# com.oraclexbt.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.oraclexbt.agent</string>
    
    <!-- Program to Run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zach/Downloads/oraclyst-python-sdk-main/run_oracle.sh</string>
    </array>
    
    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/zach/Downloads/oraclyst-python-sdk-main</string>
    
    <!-- Run on Load (at startup) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep Alive - restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <!-- Standard Output -->
    <key>StandardOutPath</key>
    <string>/Users/zach/Downloads/oraclyst-python-sdk-main/oracle_stdout.log</string>
    
    <!-- Standard Error -->
    <key>StandardErrorPath</key>
    <string>/Users/zach/Downloads/oraclyst-python-sdk-main/oracle_stderr.log</string>
    
    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <!-- Throttle Restart (don't restart too quickly) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
